# crontab
10 0 1 * *	curl -s https://harpers-wrwayefrzd.now.sh >$HOME/.motd
